# deploy/supervisor/auctionball-scheduler.conf
; Supervisor config for the Laravel scheduler.
;
; Replaces the traditional `crontab` entry. Drives:
;   - Subscription renewal sweeps
;   - Renewal-reminder emails (3 days before expiry)
;   - Audit log retention pruning (if configured)
;
; Pick ONE of (a) cron entry OR (b) this supervisor process — not both.
;
; Cron alternative (every minute, classic Laravel):
;   * * * * * cd /var/www/auctionball && php artisan schedule:run >> /dev/null 2>&1
;
; Supervisor advantage: stays warm in memory; same process model as queue/reverb;
; one tool to monitor everything.

[program:auctionball-scheduler]
process_name=%(program_name)s
command=php /var/www/auctionball/artisan schedule:work
directory=/var/www/auctionball
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/auctionball/scheduler.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
