# deploy/supervisor/auctionball-queue.conf
; Supervisor config for the AuctionBall queue worker.
;
; Without this running, queued jobs (subscription renewals, email sending,
; PDF generation, etc.) accumulate in the `jobs` table and never execute.
;
; Install on production:
;   sudo cp deploy/supervisor/auctionball-queue.conf /etc/supervisor/conf.d/
;   # edit `directory`, `user`, and `command` paths to match your server
;   sudo supervisorctl reread
;   sudo supervisorctl update
;   sudo supervisorctl start auctionball-queue:*
;
; Verify:
;   sudo supervisorctl status auctionball-queue:*
;   tail -f /var/log/auctionball/queue.log

[program:auctionball-queue]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/auctionball/artisan queue:work redis --sleep=3 --tries=3 --max-time=3600 --max-jobs=1000
directory=/var/www/auctionball
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/log/auctionball/queue.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stopwaitsecs=3600
