# deploy/supervisor/auctionball-reverb.conf
; Supervisor config for the Laravel Reverb WebSocket server.
;
; Reverb is what powers live auction broadcasts (bid placed, player changed,
; sold/unsold, timer extended). Without this running, the dashboard and
; big-screen will fall back to polling and feel laggy.
;
; In production, run Reverb behind nginx/Apache as a reverse proxy on TLS so
; clients can connect over wss://. Internal port (8080) is the bind target.
;
; Install on production:
;   sudo cp deploy/supervisor/auctionball-reverb.conf /etc/supervisor/conf.d/
;   sudo supervisorctl reread && sudo supervisorctl update
;   sudo supervisorctl start auctionball-reverb
;
; Verify:
;   sudo supervisorctl status auctionball-reverb
;   tail -f /var/log/auctionball/reverb.log
;   curl -sI http://127.0.0.1:8080  # should see Reverb headers

[program:auctionball-reverb]
process_name=%(program_name)s
command=php /var/www/auctionball/artisan reverb:start --host=0.0.0.0 --port=8080
directory=/var/www/auctionball
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/auctionball/reverb.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stopwaitsecs=10
